# supervisor.ini
[program:pgcrcollector]
command=python3 -u pgcr_collector.py
directory=/app
autostart=true
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:pgcrconsumer]
command=python3 -u pgcr_consumer.py
directory=/app
autostart=true
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:statcollector]
command=python3 -u stat_collector.py
directory=/app
autostart=true
numprocs=6
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0