# configs/supervisor.ini
[program:pgcrcollector]
command=/venv/bin/python3 -u ./app/workers/pgcr_collector.py
directory=/app
autostart=true
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:statcollector]
command=/venv/bin/python3 -u ./app/workers/roster_consumer.py
directory=/app
autostart=true
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
